# systemd/distress.service.d/itarmybox-upload-cap.conf
[Service]
ExecStartPre=/usr/bin/env php /var/www/html/itarmybox-webui/distress_service_prepare.php
